# systemd/hospital_monitor.service
[Unit]
Description=Eastern Health Hospital Wait Time Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi-zwapp
WorkingDirectory=/home/pi-zwapp/hospital-monitor
ExecStart=/home/pi-zwapp/hospital-monitor/run_monitor.sh
StandardOutput=journal
StandardError=journal
